# rpi-sky-shooter.service
[Unit]
Description= shoot the sky image.

[Service]
Type=simple
EnvironmentFile=/etc/rpi-sky-camera/environment
ExecStart=/usr/bin/python3 /etc/rpi-sky-camera/rpi-sky-shooter.py

[Install]
WantedBy=multi-user.target
